# yf_daily/cron/yf.cron
# 컨테이너 안에서 서울시간으로 실행
CRON_TZ=Asia/Seoul

# 부팅(컨테이너 시작) 즉시 1회 실행

# 매일 09:05에 1회 실행
5 9 * * * python -u /app/app.py --once
